# systemd/ferrosdr-watchdog.timer
[Unit]
Description=FerroSDR USB watchdog timer

[Timer]
OnActiveSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
