# systemd/ism-wifi-skymap3d.service
[Unit]
Description=raspi81 3D Satellite Skymap (port 8094)
After=network.target ism-wifi-gps.service

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/ism-wifi-monitor
ExecStartPre=/bin/sleep 12
ExecStart=/home/user/ism-wifi-monitor/venv/bin/python /home/user/ism-wifi-monitor/skymap3d.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
